# services/kylo/run
#!/usr/bin/env bash

sds_host=`echo ${KYLO_SPRING_DATASOURCE_URL} | cut -d/ -f3 | cut -d: -f1`
sds_port=`echo ${KYLO_SPRING_DATASOURCE_URL} | cut -d/ -f3 | cut -d: -f2`
my_tcp_wait ${sds_host} ${sds_port}

my_service "wait" ${HIVE_SERVICE_NAME}

jms_host=`echo ${KYLO_JMS_ACTIVEMQ_BROKER_URL} | cut -d/ -f3 | cut -d: -f1`
jms_port=`echo ${KYLO_JMS_ACTIVEMQ_BROKER_URL} | cut -d/ -f3 | cut -d: -f2`
my_tcp_wait ${jms_host} ${jms_port}

kylo-service start
sleep infinity &
wait $!
